# crontab/cronnightlight
0 20 * * * zsh /home/awakefox/.dotfiles/bin/.local/scripts/nightlight
