# crontab/jobs
0 15 * * * zsh /home/awakefox/.dotfiles/bin/.local/scripts/nightlight
