# crontab/jobs
50 16 * * * zsh /home/awakefox/.dotfiles/bin/.local/scripts/nightlight
